# com.acesstats.autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.acesstats.autopush</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/stephengiordano/Documents/Aces/AcesStats/autopush.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/stephengiordano/Library/Logs/acesstats-autopush.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/stephengiordano/Library/Logs/acesstats-autopush.log</string>
</dict>
</plist>
